# launch_agents/com.arca.mycloud-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.arca.mycloud-watchdog</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/danexall/biomimetics/scripts/maintain_mycloud_mount.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/danexall/biomimetics</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>StartInterval</key>
    <integer>60</integer>
    
    <key>StandardOutPath</key>
    <string>/Users/danexall/.arca/mycloud_watchdog.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/danexall/.arca/mycloud_watchdog.err</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    
    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
