# launch_agents/com.arca.omni-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.arca.omni-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.13/bin/python3</string>
        <string>/Users/danexall/omni_sync.py</string>
        <string>--config</string>
        <string>/Users/danexall/omni_sync_config.json</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/danexall</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>StandardOutPath</key>
    <string>/Users/danexall/.arca/omni_sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/danexall/.arca/omni_sync.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>PYTHONIOENCODING</key>
        <string>UTF-8</string>
    </dict>
</dict>
</plist>
